# tmp.cron
@reboot /var/www/BookwormDB/startup.sh
@reboot /var/www/BookwormDB_phonemes/startup.sh
@reboot /var/www/TAL_Pause/startup.sh
